# realtime_ml.crontab
CRON_TZ=UTC
35 13 * * * cd /home/tyreekfrazier/ISU_Research_LOCAL_RUN/mesoanalysis/gempak-scripts && ./publish_latest_ml_output.sh >> /home/tyreekfrazier/ISU_Research_LOCAL_RUN/mesoanalysis/gempak-scripts/cron_logs/realtime_ml_publish.log 2>&1
30 18 * * * cd /home/tyreekfrazier/ISU_Research_LOCAL_RUN/mesoanalysis/gempak-scripts && ./publish_verification_output.sh $(date -u -d 'yesterday' +\%Y\%m\%d) >> /home/tyreekfrazier/ISU_Research_LOCAL_RUN/mesoanalysis/gempak-scripts/cron_logs/realtime_ml_verification_publish.log 2>&1
